# slave.conf
[supervisord]
nodaemon=true
[program:spark-slave]
command=/opt/spark/sbin/start-slave.sh spark://${SPARK_MASTER_IP}:7077
